# .config/systemd/user/redshift.service.d/override.conf
[Unit]
PartOf=graphical-session.target

[Service]
Slice=session.slice
